# src/supervisord.conf
[supervisord]
nodaemon = true
environment = GRAPHITE_STORAGE_DIR='/var/lib/graphite/whisper',GRAPHITE_CONF_DIR='/etc/graphite/'

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:carbon-cache]
environment = PYTHONPATH='/var/lib/graphite/lib:$PYTHONPATH',GRAPHITE_ROOT="/usr/share/graphite-web"
user = root
command = /usr/bin/python /usr/bin/carbon-cache --debug --config=/etc/carbon/carbon.conf start
stdout_logfile = /var/log/graphite/%(program_name)s.log
stderr_logfile = /var/log/graphite/%(program_name)s.log
autorestart = true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_errfile=/var/log/supervisor/%(program_name)s.log
autorestart=true